# .config/systemd/user/spotify_player.service
[Unit]
Description=Spotify Player (user service)
Documentation=https://github.com/aome510/spotify-player
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=forking
# Leave User/Group out of user units; systemd --user already runs as the user
# WorkingDirectory can be the user's config dir
ExecStart=%h/.cargo/bin/spotify_player --daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
